# wttr.timer
[Unit]
Description=Runs wttr hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
